# config/forum.conf
[program:forum]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.config.py
directory=/root/forum
autostart=true
autorestart=true
