# config/forum.conf
[program:Forum]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.config.py
directory=/root/Forum
autostart=true
autorestart=true
